# udev-rules/etc/udev/rules.d/99-function-trackpad.rules
ACTION=="bind", SUBSYSTEM=="hid", ENV{HID_NAME}=="Function Trackpad", RUN+="/home/wade/dotfiles/scripts/.config/scripts/udev-function-trackpad.sh"
