# udev-rules/etc/udev/rules.d/99-huke-trackpad.rules
ACTION=="bind", SUBSYSTEM=="hid", ENV{HID_NAME}=="HUKE T3 Pro", RUN+="/home/wade/dotfiles/scripts/.config/scripts/udev-huke-trackpad.sh"
