# systemd/.config/systemd/user/change-wallpaper.service
[Unit]
Description=Change Wallpaper

[Service]
Type=oneshot
ExecStart=%h/.config/scripts/wallpaper.sh
