# systemd/.config/systemd/user/dwm-session.target
[Unit]
Description=DWM Graphical Session
BindsTo=graphical-session.target
Before=graphical-session.target
